# docker/rootfs/collabspace/etc/services.d/nginx/run
#!/usr/bin/bash

set -x

mkdir --parents /var/cache/nginx/cache

exec nginx -g "daemon off;"
